# lesson.6.systemd/1/val-for-lesson.timer
[Unit]
Description=timer val-for-lesson
 
[Timer]
OnCalendar=*:*:0,30
Unit=val-for-lesson.service

[Install]
WantedBy=timers.target
